# scripts/iot-web3.service
[Unit]
Description=IOT Web3 System Startup Service
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/IOT_Web3/scripts/startup-iot-web3.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target

